# ops/systemd/ollama-embed.service
[Unit]
Description=Ollama embedding model (GPU1, port 11435)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment="CUDA_VISIBLE_DEVICES=1"
Environment="CUDA_DEVICE_ORDER=PCI_BUS_ID"
Environment="OLLAMA_HOST=0.0.0.0:11435"
Environment="OLLAMA_MODELS=/var/cache/ollama"
Environment="OLLAMA_KEEP_ALIVE=-1"
Environment="OLLAMA_MAX_LOADED_MODELS=1"
ExecStart=/usr/local/bin/ollama serve
Restart=always
RestartSec=5
TimeoutStopSec=30
User=ollama
Group=ollama

[Install]
WantedBy=multi-user.target
